# package/base-file/etc/inittab
# inittab for linux, poeigl 1.3
# Format:
# ttyline:termcap-entry:getty-command
tty1:console:/etc/getty 9600 tty1
tty2:console:/etc/getty 9600 tty2
tty3:console:/etc/getty 9600 tty3
tty4:console:/etc/getty 9600 tty4
# tty5:console:/etc/getty 9600 tty5
# ttys1:dumb:/etc/getty 9600 ttys1
# ttys2:dumb:/etc/getty -m -t60 2400 ttys2
